# supervisor.conf
[program:landing-page]
command=java -jar /opt/landing-page/landing-page.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/landing-page.err.log
stdout_logfile=/var/log/landing-page.out.log
